# strfrygui.service
[Unit]
Description=StrfryGUI - Nostr Relay Management Portal
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/adam/repos/strfrygui
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="STRFRY_BINARY=/usr/local/bin/strfry"
Environment="STRFRY_CONFIG=/etc/strfry.conf"
ExecStart=/usr/bin/python3 -m flask run --host=127.0.0.1 --port=5000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
